# backend/smartdoc-ai.service
[Unit]
Description=SmartDoc AI Backend
After=network.target redis.service

[Service]
Type=simple
WorkingDirectory=/opt/smartdoc-ai/backend
EnvironmentFile=/opt/smartdoc-ai/backend/.env
ExecStart=/usr/local/bin/uv run uvicorn main:app --host 0.0.0.0 --port 8000 --access-log --log-level info
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
